# package/postinstall.sh
#!/bin/bash
#

systemctl enable broadtail
systemctl start broadtail