# package/preremove.sh
#!/bin/bash
#

systemctl stop broadtail
systemctl disable broadtail